# runit/x-ui/run
#!/bin/sh
exec 2>&1
echo -e  "starting x-ui..."
cd /usr/local/x-ui/ && exec chpst -e /usr/local/x-ui/ /usr/local/x-ui/x-ui